# attempt-1-debian/debian/post-install/etc/systemd/system/zfs-shared.target
[Unit]
Description=Shared ZFS startup target
Requires=zfs-shared.service

[Install]
WantedBy=multi-user.target
