# attempt-1-debian/debian/post-install/etc/systemd/system/zpool-shared.target
[Unit]
Description=Shared ZPOOL startup target
Requires=zpool-shared.service

[Install]
WantedBy=zfs-shared.target
